# rules.d/10-persistent-network.rules
ACTION=="add", SUBSYSTEM=="net",ATTR{address}=="00:01:02:03:04:05",NAME="usb0",RUN+="/opt/netlab/bin/addUsbNic.sh"
